# crontab
# This file needs a new line at the end
0 * * * * /usr/local/bin/pipenv run python /run.py >> /var/log/cron.log 2>&1
