# crontab
# This file needs a new line at the end
0 * * * * pipenv run python /run.py
* * * * * echo "Hello world" >> /var/log/cron.log 2>&1
